# saltmaster/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
user=root
command=/usr/sbin/sshd -D
autorestart=true

[program:salt-master]
user=root
command=/usr/bin/salt-master -l debug
autorestart=true

