# doop/etc/supervisord-zookeeper.conf
[program:zookeeper]
user=zookeeper
command=/usr/lib/zookeeper/bin/zkServer.sh start-foreground
autorestart=true
startsecs=10
startretries=3
environment=JAVA_HOME="/usr/lib/java"

