# doop/etc/supervisord-zookeeper.conf
[program:zookeeper]
user=zookeeper
command=/usr/bin/zookeeper-server start-foreground
autorestart=true
startsecs=10
startretries=3
environment=JAVA_HOME="/usr/lib/java"

