# accumulo/supervisord-accumulo-tserver.conf
[program:accumulo-tserver]
user=accumulo
command=/usr/lib/accumulo/bin/accumulo tserver
autorestart=true
startsecs=10
startretries=3
priority=250
environment=JAVA_HOME="/usr/lib/jvm/java-1.7.0"

